# scripts/cert-sync/cert-sync.service
[Unit]
Description=Daily Check for New Certificates and Deployment

[Service]
Type=oneshot
ExecStart=/opt/scripts/cert-sync.sh