# scripts/cert-sync/cert-sync.timer
[Unit]
Description=Run file cert sync script daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target